# scripts/pearl_star/systemd/postgresql@17-main.service.d/pearl-star.conf
# Pearl Star Job Queue V1 — postgresql@17-main drop-in
#
# Install to: /etc/systemd/system/postgresql@17-main.service.d/pearl-star.conf
# Then:       sudo systemctl daemon-reload && sudo systemctl restart postgresql@17-main
#
# Spec: PEARL_STAR_JOB_QUEUE_V1_SPEC.md §4.5 (ACID durability), §4.6
#       (reboot-resume: broker comes up before workers), §14 (single-NVMe
#       => the queue's broker must restart hard on failure).
#
# This drop-in does NOT change Postgres tuning (synchronous_commit is set via
# ALTER SYSTEM in 01_postgres17.sh). It only hardens restart + ordering so the
# broker is reliably up at boot before the Procrastinate workers start.
[Unit]
# Broker is the root of the durable work-list; treat its absence as critical.
Description=PostgreSQL 17 (Pearl Star Job Queue broker) — hardened restart

[Service]
# Restart the broker on any unclean exit (spec §4.6 reboot-resume substrate).
Restart=on-failure
RestartSec=5s
# Postgres can be slow to recover WAL after a power loss; give it room.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
