# scripts/pearl_star/systemd/pearl-star-monitor.service
# Pearl Star Job Queue V1 — periodic health monitor + alert emitter
#
# Install to: /etc/systemd/system/pearl-star-monitor.service
# Then:       sudo systemctl daemon-reload && sudo systemctl enable --now pearl-star-monitor
#
# Spec: PEARL_STAR_JOB_QUEUE_V1_SPEC.md §5.4 (operator file-drop alerts:
#       QUEUE_DEPTH_HIGH / VRAM_SATURATED / BROKER_HEARTBEAT_MISS). Source:
#       bin/monitor.py. This is the always-on light health sampler, distinct
#       from the per-job watchdog. Phase C upgrades it to Pearl_PM tracker /
#       Prometheus integration.
[Unit]
Description=Pearl Star Job Queue monitor (queue depth + VRAM + broker liveness alerts)
After=network-online.target postgresql@17-main.service
Wants=postgresql@17-main.service

[Service]
Type=simple
User=pearl-star
Group=pearl-star
WorkingDirectory=/opt/pearl-star/app
EnvironmentFile=/etc/pearl-star/queue.env
Environment=PYTHONPATH=/opt/pearl-star/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/pearl-star/venv/bin/python /opt/pearl-star/app/monitor.py
Restart=always
RestartSec=10s
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/var/lib/pearl-star /var/log/pearl-star

[Install]
WantedBy=multi-user.target
