# scripts/pearl_star/systemd/pearl-star-watchdog.service
# Pearl Star Job Queue V1 — stall/crash watchdog
#
# Install to: /etc/systemd/system/pearl-star-watchdog.service
# Then:       sudo systemctl daemon-reload && sudo systemctl enable --now pearl-star-watchdog
#
# Spec: PEARL_STAR_JOB_QUEUE_V1_SPEC.md §4.6 (reboot-resume unit), §5.2
#       (watchdog 60 s tick; 2x/3x thresholds; SIGTERM->grace->SIGKILL;
#       heartbeat-silent>90 s => CRASHED). Source: bin/watchdog.py.
#
# Runs the watchdog as a long-lived daemon (its own internal 60 s loop). It must
# be able to signal worker PIDs, so it runs as root (kill across users) but reads
# only heartbeat journals + nvidia-smi; it does NOT touch the broker.
[Unit]
Description=Pearl Star Job Queue watchdog (stall + crash detection)
After=network-online.target procrastinate-worker.service
Wants=procrastinate-worker.service

[Service]
Type=simple
# root so it can SIGTERM/SIGKILL pearl-star worker PIDs (spec §5.2 AUTO-KILL).
User=root
WorkingDirectory=/opt/pearl-star/app
EnvironmentFile=/etc/pearl-star/queue.env
Environment=PYTHONPATH=/opt/pearl-star/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/pearl-star/venv/bin/python /opt/pearl-star/app/watchdog.py
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
