# scripts/pearl_star/systemd/procrastinate-worker-llm.service
# Pearl Star Job Queue V1 — Procrastinate worker (llm CJK translate)
#
# Install to: /etc/systemd/system/procrastinate-worker-llm.service
# Deploy alongside t2i worker; shares gpu_heavy flock (spec §4.7).
#
# Spec: PEARL_STAR_JOB_QUEUE_V1_SPEC.md §3.3 (llm queue), OPD-20260629-003.
[Unit]
Description=Pearl Star Procrastinate worker (llm translate_atoms CJK)
After=network-online.target postgresql@17-main.service
Wants=network-online.target postgresql@17-main.service
Requires=postgresql@17-main.service

[Service]
Type=simple
User=pearl-star
Group=pearl-star
WorkingDirectory=/opt/pearl-star/app
EnvironmentFile=/etc/pearl-star/queue.env
Environment=PROCRASTINATE_APP=app.app
Environment=PYTHONPATH=/opt/pearl-star/app
Environment=PYTHONUNBUFFERED=1
Environment=PS_PHOENIX_REPO=/home/ahjan108/phoenix_omega

ExecStartPre=/bin/sh -c 'while [ -f "${PS_LIB_DIR:-/var/lib/pearl-star}/queue.paused" ]; do sleep 3; done'

# NOTE: --app uses dot-form `app.app` (module.attr), NOT colon-form `app.app:app`.
# procrastinate mis-parses the `:` and fails to import the app (Pearl Star kit #1692).
ExecStart=/opt/pearl-star/venv/bin/python -m procrastinate \
    --app app.app worker --queues llm --concurrency 1 --name llm

Restart=always
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=120

NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/var/lib/pearl-star /var/log/pearl-star /run/pearl-star /home/ahjan108/phoenix_omega

[Install]
WantedBy=multi-user.target
